# deploy/systemd/521wolf-worker.service.example
[Unit]
Description=521wolf PostgreSQL task worker
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/521wolf/app
EnvironmentFile=/opt/521wolf/.env
Environment=WOLF_USE_PG_TASK_QUEUE=true
Environment=PG_POOL_ROLE=worker
Environment=WOLF_PROCESS_ROLE=worker
ExecStart=/usr/bin/env uv run python -m app.tools.run_ui_task_worker
Restart=always
RestartSec=5

# Keep logs in journald. Inspect with:
#   journalctl -u 521wolf-worker -f

[Install]
WantedBy=multi-user.target
